# configs/systemd/grub-btrfs-refresh.timer
[Unit]
Description=Periodic grub-btrfs menu refresh fallback

[Timer]
OnBootSec=3min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
