# configs/systemd/grub-btrfs-refresh.service
[Unit]
Description=Refresh grub-btrfs snapshot menu safely

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/grub-btrfs-safe-refresh
